# reboot.service
[Unit]
Description=Daily Reboot

[Service]
ExecStart=sudo reboot
